# files/systemd/NAS-ebooks.automount
[Unit]
Description=Automount markus home folder
ConditionPathExists=/NAS/ebooks

[Automount]
Where=/NAS/ebooks
TimeoutIdleSec=30

[Install]
WantedBy=multi-user.target
