# files/systemd/NAS-music.automount
[Unit]
Description=Automount NAS music folder
ConditionPathExists=/NAS/music

[Automount]
Where=/NAS/music
TimeoutIdleSec=30

[Install]
WantedBy=multi-user.target
