# files/systemd/NAS-share.automount
[Unit]
Description=Automount markus home folder
ConditionPathExists=/NAS/share

[Automount]
Where=/NAS/share
TimeoutIdleSec=30

[Install]
WantedBy=multi-user.target
